# weblogic/oim
### BEGIN INIT INFO
# Provides:          oim
# Required-Start:    $network $remote_fs $syslog
# Should-Start:
# Required-Stop:     $network $remote_fs $syslog
# Should-Stop:
# Default-Start:     3 5
# Default-Stop:      0 1 2 6
# Short-Description: Oracle Identity Manager
# Description:       This script will the services on both nodes of the cluster.
### END INIT INFO
ENV="UAT"
if ! [ $(id -u) = 0 ];
then
    START_COMMAND=" "
else
    START_COMMAND="su - oimsvc -c "
fi
case "$1" in
    start)
        logger "Starting Oracle Identity Manager - $ENV"
        ${START_COMMAND} "/home/oimsvc/bin/startoim.sh"
        ;;
    stop)
        logger "Stopping Oracle Identity Manager - $ENV"
        ${START_COMMAND} "/home/oimsvc/bin/stopoim.sh"
        ;;
    status)
        ${START_COMMAND} "/home/oimsvc/bin/statusoim.sh"
        ;;
    *)
        echo "Usage: $0 {start|stop|status}"
        exit 1
        ;;
esac
exit 0
